# packer/files/etc/systemd/system/docker.service.d/exec.conf
[Service]
ExecStart=
ExecStart=/usr/bin/docker daemon -H fd:// --exec-opt native.cgroupdriver=cgroupfs --log-driver=syslog
