# packer/files/etc/systemd/system/docker.service.d/standalone.conf
[Service]
ExecStart=
ExecStart=/usr/bin/docker daemon
